# config/systemd/user/mesh1024.socket
[Unit]
Description=meshcentral 1024 socket
Conflicts=mesh1024.service

[Socket]
ListenStream=1024
Accept=no

[Install]
WantedBy=sockets.target
