# config/systemd/user/www8080.service
[Unit]
Description=http 8080 service
After=network.target

[Service]
ExecStart=/home/mark/go/bin/httpDir
WorkingDirectory=/opt/www8080
